# nfs-autoshare-daemon/nfs-autoshare.service
[Unit]
Description=NFS zero config daemon

[Service]
ExecStart=/usr/sbin/nfs-autoshare-daemon

[Install]
WantedBy=multi-user.target